# config/daas.conf
respawn
respawn limit 15 5

start on runlevel [2345]
stop on runlevel [06]

script
su - root -c "/root/.rbenv/shims/ruby /opt/DaaS/ws.rb 2>&1" >>  /var/log/upstart/daas.log
end script
